# supervisord.conf
# ============================================================
#  supervisord.conf — Hugging Face Spaces
#  Manages 4 processes in one container:
#    1. detection_engine  — MediaPipe + TF inference
#    2. gesture_app       — Flask dashboard, port 5000
#    3. avatar_app        — Flask avatar,    port 5001
#    4. nginx             — Reverse proxy,   port 7860
# ============================================================

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# ── 1. ISL Detection Engine (starts first — loads TF model) ──
[program:detection_engine]
command=python /app/isl_detection.py
directory=/app
autostart=true
autorestart=true
priority=1
startretries=5
startsecs=20
stdout_logfile=/tmp/detection.log
stderr_logfile=/tmp/detection_err.log
stdout_logfile_maxbytes=5MB
environment=
    RUN_MODE="CLOUD",
    PYTHONUNBUFFERED="1",
    TF_CPP_MIN_LOG_LEVEL="3",
    CUDA_VISIBLE_DEVICES="-1"

# ── 2. Gesture Dashboard (Flask + Socket.IO, port 5000) ───────
[program:gesture_app]
command=python /app/isl_ui_dashboard.py
directory=/app
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=15
stdout_logfile=/tmp/gesture.log
stderr_logfile=/tmp/gesture_err.log
stdout_logfile_maxbytes=5MB
environment=
    RUN_MODE="CLOUD",
    PYTHONUNBUFFERED="1",
    TF_CPP_MIN_LOG_LEVEL="3",
    CUDA_VISIBLE_DEVICES="-1"

# ── 3. Avatar App (Flask + eventlet, port 5001) ───────────────
[program:avatar_app]
command=python /app/speech_to_sign_avatar.py
directory=/app
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=10
stdout_logfile=/tmp/avatar.log
stderr_logfile=/tmp/avatar_err.log
stdout_logfile_maxbytes=5MB
environment=
    RUN_MODE="CLOUD",
    PYTHONUNBUFFERED="1",
    PORT="5001"

# ── 4. nginx (starts after Flask apps are up) ─────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /etc/nginx/nginx.conf
autostart=true
autorestart=true
priority=30
startretries=3
startsecs=5
stdout_logfile=/tmp/nginx.log
stderr_logfile=/tmp/nginx_err.log
stdout_logfile_maxbytes=5MB
